# servers/sanat/roles/mediawiki/templates/mw-jobrunner.service.j2
[Unit]
Description=MediaWiki jobrunner for {{ instance.key }}

[Service]
ExecStart=/srv/mediawiki/{{ instance.key }}/targets/production/maintenance/run runJobs --wait --maxjobs=1000 --memory-limit=2G --procs=2
Restart=always
RestartSec=5
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
Nice=15

[Install]
WantedBy=multi-user.target
