# servers/sanat/roles/mediawiki/templates/mw-jobrunner.service.j2
[Unit]
Description=MediaWiki jobrunner for {{ instance }}
After=network.target
Requires=mariadb.service

[Service]
user=hhvm
Environment=HHVM_REPO_CENTRAL_PATH=/run/hhvm/hhvm.hhbc
ExecStart=/usr/bin/php /srv/mediawiki/targets/production/maintenance/runJobs.php --wait --maxjobs=10000 --wiki={{ instance }}
Restart=always
RestartSec=5
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelTunables=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
